# ubuntu_helpers/rn_irl.service
[Unit]

Description=Really Nice IRL
After=network.target
StartLimitIntervalSec=0

[Service]

Type=simple
Restart=always
RestartSec=1
User=root
ExecStartPre=
ExecStart=/bin/rn_irl
ExecStartPost=
ExecStop=
ExecReload=

[Install]

WantedBy=multi-user.target
